# .config/systemd/user/hyprland.service
[Unit]
Description=Start Hyprland with systemd
BindsTo=hyprland.target
PartOf=hyprland.target
Before=hyprland.target
After=systemd-user-sessions.target

[Service]
Type=simple
Environment=LIBSEAT_BACKEND=logind
ExecStart=/bin/bash -c '/usr/bin/Hyprland'
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
