# .config/systemd/user/hyprland.target
[Unit]
Description=Hyprland Target
Documentation=man:systemd.special(7)
Requires=hyprland.service
After=hyprland.service
AllowIsolate=yes
